# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/data/cfg/video_default_cfg.ini
/sbin/sound
/data/cfg/sensor_db.cfg
/data/cfg/buildroot_cfg.ini
/etc/resolv.conf
/sbin/factorywrite
/etc/init.d/done
/var
/lib/firstboot/99_10_no_fo_cleanup
/etc/hotplug2-init.rules
/bin/resetConfigReload.sh
/etc/preinit
/lib/firstboot/50_pivot
/sbin/mount_root
/lib/firstboot/30_reset_copy_rom
/etc/hotplug2-common.rules
/etc/shells
/data/cfg/audio_default_cfg.ini
/lib/preinit/50_choose_console
/lib/firstboot/10_no_fo_clear_overlay
/etc/protocols
/lib/preinit/20_check_jffs2_ready
/data/cfg/loadSensor
/lib/preinit/10_indicate_preinit
/etc/config/protocol
/etc/init.d/check_dup_list_from_update
/lib/firstboot/99_10_with_fo_cleanup
/etc/hotplug.d/button/00-button
/etc/config/system
/bin/timezone.sh
/etc/fstab
/lib/preinit/30_device_fs_daemons
/lib/preinit/90_mount_no_jffs2
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/lib/firstboot/20_has_mini_fo
/lib/firstboot/10_reset_has_mini_fo
/lib/preinit/40_dont_mount_jffs2
/lib/firstboot/10_determine_parts
/sbin/led.sh
/lib/firstboot/40_no_fo_copy_ramoverlay
/etc/banner
/lib/preinit/10_essential_fs
/etc/shadow
/etc/rc.common
/lib/preinit/40_init_shm
/init
/bin/resetWaitUnbind.sh
/lib/preinit/20_device_fs_mount
/lib/preinit/80_mount_root
/etc/init.d/sysntpd
/etc/rc.local
/lib/functions/uci-defaults.sh
/lib/functions/parse_uci_config.sh
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/firstboot/30_is_rootfs_mounted
/sbin/hotplug-call
/etc/oem_id.bin
/lib/preinit/40_mount_devpts
/lib/preinit/90_init_console
/etc/services
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/10_check_for_mtd
/lib/upgrade/common.sh
/sbin/core_helper
/etc/sysupgrade.conf
/lib/firstboot/05_firstboot_skip
/lib/upgrade/keep.d/base-files
/bin/stat_check.sh
/bin/onboarding.sh
/etc/group
/bin/resetWifiOnly.sh
/etc/functions.sh
/bin/logControl.sh
/sbin/firstboot
/bin/dsdSendCmd.sh
/bin/login.sh
/data/cfg/isp_default_cfg.ini
/etc/hosts
/etc/hotplug.d/usb/10-usb
/lib/functions.sh
/lib/preinit/60_init_hotplug
/etc/config/network
/etc/init.d/rcS
/rom/note
/etc/inittab
/etc/openwrt_release
/lib/preinit/05_mount_skip
/sbin/wifi
/etc/profile
/etc/init.d/check_upgrade
/etc/init.d/sysfixtime
/etc/TZ
/lib/firstboot/20_no_fo_mount_jffs
/etc/hotplug.d/net/09-usb-if
/lib/functions/service.sh
/etc/init.d/umount
/lib/upgrade/keep.d/base-files-essential
/lib/firstboot/30_no_fo_pivot
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/etc/mtab
/lib/firstboot/40_copy_ramoverlay
/etc/init.d/boot
/lib/functions/boot.sh
/bin/ipcalc.sh
/etc/diag.sh
/bin/tapo_kill.sh
/lib/firstboot/20_reset_clear_jffs
