# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/diagnose.list
/usr/bin/consoleredict
/usr/bin/logrecordd
/etc/init.d/logrecord
