# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/dsd.list
/etc/config/luci
/etc/config/cover
/etc/config/smart_detection
/etc/config/unusual_detection
/etc/config/msg_alarm
/etc/config/lens_mask
/etc/config/smart_data
/etc/config/abnormal_events
/sbin/luci-reload
/usr/lib/libdsd_client.so
/etc/config/audio_capability
/sbin/debug_ir_cut_switch_test.sh
/etc/config/tamper_detection
/etc/config/ucitrack
/etc/config/ai_enhance_capability
/etc/config/msg_alarm_plan
/etc/config/people_detection
/sbin/track-reload
/etc/config/intrusion_detection
/etc/config/target_track
/usr/bin/dsd_exec
/etc/config/smart_msg_push_capability
/etc/init.d/dsd
/lib/upgrade/luci-add-conffiles.sh
/sbin/debug_ptz_test.sh
/etc/config/user_management
/etc/config/preset
/etc/config/motion_detection
/etc/config/pattern
/etc/config/motion_detection_cloud
/etc/config/timing_reboot
/etc/config/OSD
/sbin/debug_speaker_test.sh
/etc/config/playback
/etc/config/ROI
/etc/dsd_convert.json
/usr/bin/dsd
/etc/config/linkage_capability
/etc/config/ptz_plan
/etc/config/linecrossing_detection
/etc/config/ptz
/etc/config/multicast
/etc/config/sound_detection
/etc/config/tour
