# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/ledd.list
/lib/modules/3.10.14/gpio_leds.ko
/etc/config/led
/usr/bin/ledd
/etc/init.d/ledd
