# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/recover-manager.list
/etc/init.d/recover-manager
/usr/bin/get_sys_status
/etc/init.d/auto_upgrade
/usr/bin/cloud_upgrade
/etc/config/auto_upgrade
/usr/bin/auto_upgrade
/usr/bin/get_new_fwinfo
/usr/bin/recover-manager
