# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/t_relayd.list
/etc/init.d/relayd
/usr/bin/relayd
/etc/config/relayd_config
/etc/root.cer
