# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/telemetry.list
/etc/init.d/telemetry
/bin/telemetry
/etc/config/telemetry
/bin/create_telemetry_files
