# firmware/copie_firmware/rom/mnt/usr/lib/opkg/info/wlan-manager.list
/usr/sbin/wlan-manager
/etc/init.d/wlan
/lib/network/is_cal.script
/lib/network/dhcp.script
/etc/init.d/wlan-manager
/etc/config/wlan
/lib/network/hotplug.script
/etc/default/hostapd.accept
/usr/sbin/collect_wifi_info
/etc/default/hostapd.deny
/etc/config/wlan_runtime
/usr/sbin/renew_ip.sh
/etc/default/hostapd_default.conf
