# firmware/copie_firmware/usr/lib/opkg/info/cet.list
/lib/audio/register_success.g711
/bin/logo1080P.bmp.bin
/bin/logo2304x1296.bmp.bin
/lib/audio/wlan_state_connect_suc.g711
/lib/audio/prompt_tone_sucess.g711
/etc/config/system_state_audio
/lib/audio/reset_success.g711
/lib/audio/wlan_state_password_err.g711
/lib/audio/motion_deteck_cue.g711
/bin/time2560x1440.bmp.bin
/lib/audio/wlan_state_mode_switch_suc.g711
/lib/audio/prompt_sd_card.g711
/lib/audio/sd_abnormal.g711
/etc/config/record_plan
/lib/audio/prompt_tone_failed.g711
/bin/time360P.bmp.bin
/etc/config/auto_reboot
/bin/time720P.bmp.bin
/bin/time2304x1296.bmp.bin
/etc/config/cet
/bin/logo360P.bmp.bin
/lib/audio/wlan_state_connecting.g711
/etc/config/record_control
/bin/logo2560x1440.bmp.bin
/etc/config/video_message
/etc/config/plan_advance
/bin/logo720P.bmp.bin
/bin/time1080P.bmp.bin
/bin/cet
/lib/audio/motion_deteck_alarm.g711
/lib/audio/wlan_state_get_ip_fail.g711
/etc/init.d/cet
