# firmware/copie_firmware/usr/lib/opkg/info/cloud-service.list
/bin/cloud-service
/etc/config/msg_push
/etc/config/msg_push_plan
/etc/init.d/cloud_service
/etc/config/ffs
/etc/cloud_service/ipc_service.cer
