# firmware/copie_firmware/usr/lib/opkg/info/event_center.list
/usr/bin/event_networkstatus
/usr/lib/libevent_center.so
/usr/bin/event_checkipconflict
/etc/netlink.d/wan_status_notify/10-network_abnormal
/etc/init.d/network_status
