# firmware/copie_firmware/usr/lib/opkg/info/kmod-k_hi_reset.list
/etc/hotplug.d/button/10-reset
/etc/init.d/gpio_reset
/lib/modules/3.10.14/k_hi_reset.ko
