# firmware/copie_firmware/usr/lib/opkg/info/t_rtspd.list
/etc/config/rtspd_config
/usr/bin/rtspd
/usr/lib/libcommon.so
/etc/init.d/rtspd
