# firmware/copie_firmware/usr/lib/opkg/info/wtd.list
/etc/config/wtd
/bin/wtd
/bin/disable_hw_wtd
/etc/init.d/wtd
