# firmware/copie_firmware/bin/etc/init.d/boot
#!/bin/sh /etc/rc.common
# Copyright (c) 2014-2015 TP-LINK Technologies CO.,LTD.
# Copyright (c) 2013 Qualcomm Atheros, Inc.
# Copyright (C) 2006-2011 OpenWrt.org

START=05

default_dst="1"
default_timing_mode="ntp"

system_config() {
	local cfg="$1"

	local hostname conloglevel timezone

	config_get hostname "$cfg" hostname 'SLP'
	echo "$hostname" > /proc/sys/kernel/hostname

	config_get conloglevel "$cfg" conloglevel
	config_get buffersize "$cfg" buffersize
	[ -z "$conloglevel" -a -z "$buffersize" ] || dmesg ${conloglevel:+-n $conloglevel} ${buffersize:+-s $buffersize}

	config_get timezone "$cfg" timezone 'UTC'
	echo "$timezone" > /tmp/TZ

	config_get zonename "$cfg" zonename
	[ -n "$zonename" ] && [ -f "/usr/share/zoneinfo/$zonename" ] && ln -s "/usr/share/zoneinfo/$zonename" /tmp/localtime

	# apply timezone to kernel
	date -k

	# config vm cache
	echo 256 > /proc/sys/vm/min_free_kbytes
	echo 5 > /proc/sys/vm/dirty_ratio

	grep -q Ingenic /proc/cpuinfo && {	
		echo 598 > /proc/sys/vm/min_free_kbytes
		echo 1 > /proc/sys/vm/overcommit_memory
		echo 50 > /proc/sys/vm/overcommit_ratio
		echo 5 > /proc/sys/vm/dirty_ratio
		echo 2 > /proc/sys/vm/dirty_background_ratio
		echo 500 > /proc/sys/vm/dirty_expire_centisecs
		echo 100 > /proc/sys/vm/dirty_writeback_centisecs
		echo 0 > /proc/sys/vm/admin_reserve_kbytes
		
		echo 16777216 > /sys/block/zram0/disksize
		mkswap /dev/zram0
		swapon /dev/zram0
		echo 100 > /proc/sys/vm/swappiness
	}
}

# Set some defaults at the first time the device boots.
# This is a light version of 'uci-defaults'. by lizheng.
# NOTE:
# All scripts in the folder /etc/uci-defaults/ are automatically executed,
# and if they exited with code 0 deleted afterwards (scripts that did not
# exit with code 0 are not deleted and will be re-executed during the next
# boot until they also successfully exit).
# You can put your own default files in that folder.
apply_uci_config() {
	cd /etc/uci-defaults || return 0
	files="$(ls)"
	[ -z "$files" ] && return 0
	mkdir -p /tmp/.uci
	for file in $files; do
		( . "./$(basename $file)" ) #&& rm -f "$file"
	done

	# Commit all sections here
	uci commit
}

set_system_config() {
	config_load system
	config_get dst_enabled dst enabled
	config_get ntp_mode basic timing_mode
	if [ "$dst_enabled" != "$default_dst" ]; then
		uci_set system dst enabled 1
		echo "set dst enable to 1"
	fi

	if [ "$ntp_mode" != "$default_timing_mode" ]; then
		uci_set system basic timing_mode "$default_timing_mode"
		echo "set timing_mode to ntp"
	fi
	uci commit system
}

start() {
	[ -f /proc/mounts ] || /sbin/mount_root
	[ -f /proc/jffs2_bbc ] && echo "S" > /proc/jffs2_bbc
	[ -f /proc/net/vlan/config ] && vconfig set_name_type DEV_PLUS_VID_NO_PAD

	mkdir -p /var/run
	mkdir -p /var/log
	mkdir -p /var/lock
	mkdir -p /tmp/usbdisk
	touch /var/log/wtmp
	touch /var/log/lastlog
	touch /tmp/resolv.conf.auto
	ln -sf /tmp/resolv.conf.auto /tmp/resolv.conf
	grep -q debugfs /proc/filesystems && mount -t debugfs debugfs /sys/kernel/debug

	load_modules /etc/modules.d/*

	apply_uci_config

	set_system_config

	config_load system
	config_foreach system_config system

	# update timezone
	timezone.sh
	
	killall -q hotplug2
	[ -x /sbin/hotplug2 ] && /sbin/hotplug2 --override --persistent \
			--set-rules-file /etc/hotplug2.rules \
			--set-coldplug-cmd /sbin/udevtrigger \
			--max-children 1 >/dev/null 2>&1 &

	# create /dev/root if it doesn't exist
	[ -e /dev/root -o -h /dev/root ] || {
		rootdev=$(awk 'BEGIN { RS=" "; FS="="; } $1 == "root" { print $2 }' < /proc/cmdline)
		[ -n "$rootdev" ] && ln -s "$rootdev" /dev/root
	}

	[ -x /bin/create_telemetry_files ] && /bin/create_telemetry_files
}

