# firmware/copie_firmware/bin/etc/init.d/nvid
#!/bin/sh /etc/rc.common

START=77
RESUME=77

start() {
    /bin/nvid 2>&1 1>/dev/console &
}
