# firmware/copie_firmware/bin/etc/init.d/recover-manager
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=98

USE_PROCD=1

start() {
	recover-manager 2>&1 1>/dev/console &
}

stop() {
	killall recover-manager
}
