# firmware/copie_firmware/rom/mnt/etc/init.d/p2pd
#!/bin/sh /etc/rc.common
#Qintian Guo <guoqintian@tp-link.net>

START=60
RESUME=60

P2PD_BIN="/usr/bin/p2pd"

start(){
	$P2PD_BIN &
}

terminate(){
	service_stop $P2PD_BIN
}

resume(){
	start
}