# firmware/copie_firmware/rom/etc/init.d/dsd
#!/bin/sh /etc/rc.common

START=30

start() {
	/usr/bin/dsd &
}

stop() {
	kill -9 `pgrep /usr/bin/dsd`
}
